# krunner/.config/autostart/krunner.desktop
[Desktop Entry]
Exec=krunner -d
Icon=
Name=krunner
Path=
Terminal=True
Type=Application
